# .config/systemd/user/switch-screen.service
[Unit]
Description=Watchdog service for display events

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rescrs
