# .config/systemd/user/switch-keymaps.service
[Unit]
Description=Watchdog service for keyboard events

[Service]
Type=oneshot
ExecStart=%h/.local/bin/remaps
